# scripts/init/helpers/99-usb-connect.rules
ACTION=="add", SUBSYSTEM=="usb", RUN+="/usr/local/bin/usb-trigger.sh"
